# front/lib/api/sandbox/image/egress/dust-egress-nftables.service
[Unit]
Description=Dust egress nftables rules for agent-proxied
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/etc/dust/egress-nftables.sh

[Install]
WantedBy=multi-user.target
