# front/lib/api/sandbox/image/egress/dust-egress-iptables.service
[Unit]
Description=Dust egress iptables rules for agent-proxied
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/etc/dust/egress-iptables.sh

[Install]
WantedBy=multi-user.target
